# deploy/ubuntu/layerpilot-backup.timer
[Unit]
Description=Run 3DSTUXXX backup every night

[Timer]
OnCalendar=*-*-* 02:30:00
Persistent=true
RandomizedDelaySec=15m
Unit=layerpilot-backup.service

[Install]
WantedBy=timers.target
