# deploy/ubuntu/layerpilot-ops-check.timer
[Unit]
Description=Run 3DSTUXXX operations health check every 15 minutes

[Timer]
OnBootSec=5m
OnUnitActiveSec=15m
Persistent=true
RandomizedDelaySec=2m
Unit=layerpilot-ops-check.service

[Install]
WantedBy=timers.target
