# deploy/ubuntu/layerpilot-ops-check.service
[Unit]
Description=3DSTU FarmFlow operations health check
Requires=docker.service
Wants=docker.service
After=docker.service

[Service]
Type=oneshot
WorkingDirectory=/opt/layerpilot
Environment=LAYERPILOT_PUBLIC_URL=http://127.0.0.1:8797
Environment=LAYERPILOT_BACKUP_DIR=/var/backups/layerpilot
Environment=LAYERPILOT_MIN_FREE_MB=2048
ExecStart=/usr/bin/env bash /opt/layerpilot/scripts/ubuntu-ops-check.sh
TimeoutStartSec=5min
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
